# files/etc/systemd/system/zeropoint-phase-cleanup.target
[Unit]
Description=Zeropoint Cleanup Services
Documentation=https://docs.zeropoint.example.com/boot-phases
After=zeropoint-agent.target
Requires=zeropoint-agent.target

[Install]
WantedBy=multi-user.target