# files/etc/systemd/system/zeropoint-setup-nvidia-post-reboot.service
[Unit]
Description=Zeropoint NVIDIA GPU Post-Reboot Verification
Documentation=https://docs.zeropoint.example.com/boot-phases
DefaultDependencies=no
After=zeropoint-setup-nvidia-drivers.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zeropoint-setup-nvidia-post-reboot.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutSec=300
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
