# files/etc/systemd/system/zeropoint-setup-storage.service
[Unit]
Description=Zeropoint First-Boot Storage Setup
Documentation=https://docs.zeropoint.example.com/boot-phases
DefaultDependencies=no
After=zeropoint-resize-rootfs.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zeropoint-setup-storage.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutSec=300
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
