# files/etc/systemd/system/zeropoint-setup-nvidia-drivers.service
[Unit]
Description=Zeropoint NVIDIA GPU Driver Installation
Documentation=https://docs.zeropoint.example.com/boot-phases
DefaultDependencies=no
After=zeropoint-update-agent.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zeropoint-setup-nvidia-drivers.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutSec=600
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
