# files/etc/systemd/system/zeropoint-resize-rootfs.service
[Unit]
Description=Zeropoint First-Boot Root Filesystem Expansion
Documentation=https://docs.zeropoint.example.com/boot-phases
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zeropoint-resize-rootfs.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutSec=300
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
